# sound/connect-sound.service
[Service]
Type=oneshot
Slice=session.slice
ExecStart=%h/dotfiles/sound/connect-sound
