# sway-session/gammastep.service
[Unit]
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/gammastep
