# sway-session/mako.service
[Unit]
PartOf=graphical-session.service

[Service]
ExecStart=/usr/bin/mako
