# sway-session/nm-applet.service
[Unit]
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/nm-applet --indicator

[Install]
WantedBy=sway-session.target
