# sway-session/polkit-gnome.service
[Unit]
PartOf=graphical-session.target

[Service]
ExecStart=/usr/lib/polkit-gnome/polkit-gnome-authentication-agent-1
