# sway-session/redshift.service
[Unit]
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/redshift -b 1.0:0.8 -t 6000:3000
