# battery-monitor.service
[Unit]
Description=lumeh Battery Monitor

[Service]
ExecStart=/home/lumpy/dotfiles/battery-monitor.sh

[Install]
