# i3-wait.service
[Service]
ExecStart=%h/bin/i3-wait
TimeoutStartSec=10s
Type=oneshot
