# mouse-monitor.service
[Unit]
After=i3-wait.service
Wants=i3-wait.service

[Service]
ExecStart=%h/bin/mouse-monitor

[Install]
WantedBy=default.target
