# mouse-monitor.service
[Service]
ExecStart=%h/bin/mouse-monitor

[Install]
WantedBy=default.target
