# update-wallpaper.timer
[Timer]
OnCalendar=*:0/30

[Install]
WantedBy=timers.target
