# update-wallpaper.service
[Service]
ExecStart=%h/bin/update-wallpaper
Type=oneshot
