# update-wallpaper.timer
[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
