# update-wallpaper.timer
[Timer]
AccuracySec=1m
OnBootSec=30m
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
