# units/gammastep.service
[Unit]
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/gammastep

[Install]
WantedBy=graphical-session.target
