# units/mako.service
[Unit]
BindsTo=compositor.service
After=compositor.service

[Service]
ExecStart=/usr/bin/mako
Restart=on-failure
Slice=session.slice
MemoryMax=2G
