# units/neovide-launcher.socket
[Socket]
ListenStream=%t/neovide-launcher.sock
Accept=yes

[Install]
WantedBy=graphical-session.target
