# units/neovide-remote@.service
[Unit]
BindsTo=graphical-session.target

[Service]
Type=exec
ExecStart=%h/bin/neovide-remote %I
Slice=app.slice
