# units/update-public-ip.service
[Service]
Type=oneshot
ExecStart=%h/bin/update-public-ip
TimeoutStartSec=30s
Slice=background.slice
