# units/update-wallpaper.service
[Unit]
StartLimitIntervalSec=0

[Service]
ExecStart=%h/bin/update-wallpaper
Type=oneshot
Slice=background.slice
