# units/update-wallpaper.timer
[Unit]
BindsTo=compositor.service
After=compositor.service

[Timer]
AccuracySec=5m
OnActiveSec=90m
OnUnitActiveSec=2h
RandomizedDelaySec=30m
