# units/wayvnc.service
[Unit]
Requires=sway-headless.service
After=sway-headless.service

[Service]
ExecStart=/usr/bin/wayvnc --gpu
TimeoutStopSec=10
