# rather-dashing/rather-dashing.service
[Unit]
Description=Rather Dashing

[Service]
Type=simple
ExecStart=/usr/local/bin/rather-dashing.py
User=colin
