# fruitron/fruitron-music-scheduler.service
[Unit]
Description=Fruitron Music Scheduler

[Service]
Type=simple
ExecStart=/usr/local/bin/fruitron-music-scheduler
User=colin
